# tftp2http.service
[Service]
WorkingDirectory=%Y
ExecStart=go run .
Environment=HOME=%Y

[Install]
WantedBy=multi-user.target
